# scripts/systemd/script-to-video-worker.service
[Unit]
Description=Script-to-Video Worker
After=network.target

[Service]
# Replace the user below with the user that owns /opt/script-to-video
User=ubuntu
Type=simple
WorkingDirectory=/opt/script-to-video
EnvironmentFile=/etc/script-to-video/worker.env
ExecStart=/usr/bin/node worker/worker.cjs
Restart=on-failure
RestartSec=5
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
